# misc/crontab
0 6     * * *     www-data cd /opt/gdprshare && ./gdprshare -cleanup > /dev/null 2>&1
